# aplication/automatizacion/sinc_pionono.service
[Unit]
Description= sinc piononoApp

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 -m /var/www/webappPIO/aplication/scripts.automaticsinc
