# aplication/automatizacion/sinc_pionono.timer
[Unit]
Description= sinc piononoApp

[Timer]
OnCalendar=*-*-* *:0/20
OnStartupSec=60
Persistent=True

[Install]
WantedBy=timers.target
