# aplication/automatizacion/sinc_pionono.service
[Unit]
Description= sinc piononoApp

[Service]
Type=oneshot
ExecStart=/var/www/webappPIO/aplication/sinc.sh
